# aio/supervisord.conf

[supervisord]
nodaemon=true
user=root

[program:postgresql]
command=/usr/lib/postgresql/13/bin/postgres -D /var/lib/postgresql/data
user=postgres
autorestart=true

[program:backend]
command=node /app/backend/index.js
directory=/app/backend
environment=PORT=4000,DB_HOST=127.0.0.1,DB_USER=leakcheck,DB_PASSWORD=secure_password,DB_NAME=leakcheck_db,JWT_SECRET=super_secret_jwt_key
user=root
autorestart=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
user=root
autorestart=true
